# dell-latitude-5411/remap-keys.service
[Unit]
Description=Remaps Dell Latitude 5411 keys

[Service]
User=michael
ExecStartPre=/bin/sleep 60
ExecStart=/usr/bin/xkbcomp /home/michael/.xkbmap ":0"

[Install]
WantedBy=multi-user.target
